# ops/vps/community-provision-operator/systemd/pirate-community-provision-operator.service
[Unit]
Description=Pirate community provision operator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pirate-spaces
Group=pirate-spaces
WorkingDirectory=/srv/pirate-spaces/app
Environment=HOME=/srv/pirate-spaces
EnvironmentFile=/srv/pirate-spaces/config/community-provision-operator.env
ExecStart=/srv/pirate-spaces/bin/start-community-provision-operator.sh
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/srv/pirate-spaces/data /srv/pirate-spaces/log /srv/pirate-spaces/run

[Install]
WantedBy=multi-user.target
